# backend/utils/podwise-batch-processor.service
[Unit]
Description=Podwise Batch Processor
After=network.target
Wants=network.target

[Service]
Type=simple
User=bai
WorkingDirectory=/home/bai/Desktop/Podwise
Environment=PATH=/home/bai/Desktop/Podwise/.venv/bin
ExecStart=/home/bai/Desktop/Podwise/.venv/bin/python backend/utils/robust_batch_processor.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# 資源限制
MemoryMax=8G
CPUQuota=80%

# 安全設定
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/bai/Desktop/Podwise

[Install]
WantedBy=multi-user.target 